# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:stf]
command=stf --help  ; replace with your desired stf command
autostart=true
autorestart=false
startsecs=0
